# plamo/05_ext/mlocate/mlocate.cron
#!/bin/sh
renice +19 -p $$ > /dev/null 2>&1
updatedb